# sync/launchd/com.budget.sync.rogers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.budget.sync.rogers</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/diogolopes/dev/budget/nextjs-boilerplate/sync/run-rogers.sh</string>
  </array>

  <!-- Fire once a day at 11:59. If the Mac is OFF at that time the run is simply
       skipped (no catch-up). If it's ASLEEP, launchd runs it on the next wake —
       which is harmless because re-imports dedup to zero. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>11</integer>
    <key>Minute</key><integer>59</integer>
  </dict>

  <!-- Don't run on load/login, only on the schedule. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/diogolopes/Library/Application Support/budget-sync/logs/rogers.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/diogolopes/Library/Application Support/budget-sync/logs/rogers.log</string>
</dict>
</plist>
