# sync/launchd/com.budget.sync.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.budget.sync.digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/diogolopes/dev/budget/nextjs-boilerplate/sync/digest.sh</string>
  </array>

  <!-- Fire at 12:30 as a backstop — each successful sync also triggers the digest
       check, so the notification normally fires the moment the last of the 4
       accounts (Rogers, Amex, Scotia, Tangerine) completes. The digest only pushes
       when all 4 have a fresh "ok" sync; otherwise it stays quiet. No browser.
       Mac OFF at 12:30 → skipped; ASLEEP → runs on next wake (harmless). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <!-- Don't run on load/login, only on the schedule. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/diogolopes/Library/Application Support/budget-sync/logs/digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/diogolopes/Library/Application Support/budget-sync/logs/digest.log</string>
</dict>
</plist>
